# blackbeard-os/os/2404-full/releng/airootfs/etc/systemd/system/getty@tty1.service.d/autologin.conf
[Service]
ExecStart=
ExecStart=-/sbin/agetty --autologin pirate --noclear %I 38400 linux
