# blackbeard-os/os/2404-full/blackbeard/airootfs/etc/systemd/system/blackbeard_installer.service
[Unit]
Description=CLI installer for blackbeard os
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/usr/local/bin/blackbeard_installer
StandardOutput=journal

[Install]
WantedBy=multi-user.target
